# extra/upstart/nero.conf
start on runlevel [2345]

respawn

exec nero-start
pre-stop exec nero-stop
